# dist/aethertftp.service.in
[Unit]
Description=AetherTFTP Headless TFTP Server (v@PROJECT_VERSION@)
After=network.target

[Service]
Type=simple
User=nobody
Group=nogroup
AmbientCapabilities=CAP_NET_BIND_SERVICE
ExecStart=/usr/bin/aethertftp --server --port 69 --dir /var/lib/tftpboot
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
